# rootfs/etc/s6-overlay/s6-rc.d/torgo-app/run
#!/command/with-contenv bash
# S6 run script for the main torgo Go application
echo "--- S6 Service: Starting torgo-app ---"
exec /app/torgo-app
