# rootfs/etc/s6-overlay/s6-rc.d/privoxy/run
#!/command/with-contenv bash
set -e
echo "--- S6 Service: Starting Privoxy ---"
# Configure Privoxy just before starting it
PRIVOXY_CONFIG_FILE="/etc/privoxy/config"

# --- CORRECTED & MINIMAL Privoxy Configuration ---
# This configuration is compatible with the version of Privoxy from the
# Alpine package repository. It removes the problematic action files.
cat > ${PRIVOXY_CONFIG_FILE} <<EOP
# --- Core Settings ---
listen-address  [::]:8118
listen-address  0.0.0.0:8118
confdir /etc/privoxy
logdir /var/log/privoxy

# --- Forwarding ---
forward-socks5t / torgo:9000 .

# --- Logging ---
# Log level is set by environment variable. 0 is quiet.
debug ${PRIVOXY_LOG_LEVEL:-0}
EOP

# Use exec to hand control over to privoxy
exec privoxy --no-daemon "${PRIVOXY_CONFIG_FILE}"
